# .config/autostart/conky.desktop
[Desktop Entry]
Name=conky
Exec=/usr/local/bin/hotplug.sh 
Type=Application
Categories=Utility
