# .config/autostart/drive_fuse.desktop
[Desktop Entry]
Name=Google Drive Fuse
Exec=/usr/bin/google-drive-ocamlfuse /mnt/gdrive 
Type=Application
Categories=Utility
