# etc/udev/rules.d/95-monitor-hotplug.rules
KERNEL=="card0", SUBSYSTEM=="drm", RUN+="/usr/bin/systemctl restart hotplug.service"

