# etc/udev/rules.d/95-monitor-hotplug2.rules
KERNEL=="card0", SUBSYSTEM=="drm", RUN+="/usr/libexec/udev_synergyc"

